# config/containers/syncthing.container
[Unit]
Description=Syncthing service

[Container]
HostName=syncthing
Image=ghcr.io/michamos/syncthing:main
Volume=syncthing.volume:/var/syncthing
